# debian/services/szarp-server.target
[Unit]
Description=SZARP Scada Software Server
After=network.target
[Install]
WantedBy=multi-user.target
