# resources/iks-server.conf
[program:iks-server]
autostart=false
autorestart=true
command=/opt/szarp/bin/iks-server
startsecs=0
