# debian/szarp-server.prerm
#!/bin/sh

systemctl stop szarp-server.target
systemctl disable szarp-server.target
