# debian/szarp-iks-server.prerm
#!/bin/sh

systemctl stop iks-server
systemctl disable iks-server
