# debian/szarp-iks-server.postinst
#!/bin/sh

set -e;

. /usr/share/debconf/confmodule;

UNITSELECT_PATH="/etc/szarp/units-enabled"

if [ ! -d $UNITSELECT_PATH ]; then
	mkdir -p $UNITSELECT_PATH
fi

touch "$UNITSELECT_PATH/iks-server"

systemctl unmask iks-server
systemctl enable iks-server
# iks-server will be started by szarp-server
